# infrastructure/mac/com.downstream.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.downstream.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/ferenczcsuszner/Coding/Experiments/AIrunBusiness && git pull origin main 2>&1 | logger -t downstream-sync</string>
    </array>

    <!-- Run every hour -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- Also run at load (startup/login) -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/downstream-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/downstream-sync.log</string>
</dict>
</plist>
